# config/systemd/user/pulseaudio.service
[Service]
ExecStart=/usr/bin/pulseaudio --daemonize=no -vvvv --log-time=1
